# deploy/system/systemd.paradar-logs.service
[Unit]
Description=Paradar log collection

[Service]
Type=oneshot
ExecStart=/opt/paradar/paradar/dump-logs.sh
SyslogIdentifier=paradar-logs
